# scripts/deployment/systemd/django.service
[Unit]
Description=Bifrost Django Admin Interface
After=network.target postgresql.service

[Service]
Type=simple
User=bifrost
WorkingDirectory=/opt/bifrost-trader-option/django_app
Environment="PATH=/opt/bifrost-trader-option/venv/bin"
ExecStart=/opt/bifrost-trader-option/venv/bin/gunicorn django_config.wsgi:application --bind 0.0.0.0:8001 --workers 2
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

